# .claude-coord/coordination-daemon.service
[Unit]
Description=Coordination Daemon for %I
After=network.target

[Service]
Type=forking
User=%i
WorkingDirectory=%h/meta-autonomous-framework
ExecStart=/usr/bin/python3 -m coord_service.daemon start --project-root %h/meta-autonomous-framework
ExecStop=/usr/bin/python3 -m coord_service.daemon stop --project-root %h/meta-autonomous-framework
PIDFile=%h/meta-autonomous-framework/.claude-coord/daemon.pid
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=append:/tmp/coord-daemon.log
StandardError=append:/tmp/coord-daemon.log

[Install]
WantedBy=default.target
